# etc/systemd/system/docker.service.d/openmediavault.conf
[Service]
EnvironmentFile=-/etc/default/docker
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// $DOCKER_OPTS \
    $OMVDOCKER_API \
    $OMVDOCKER_IMAGE_PATH
